# ops/tm-backup.timer
[Unit]
Description=Daily terminus-mind backup

[Timer]
OnCalendar=*-*-* 04:23
Persistent=true

[Install]
WantedBy=timers.target
