# ops/tm-embed.service
[Unit]
Description=llama.cpp embedding server (nomic-embed-text-v1.5) for terminus-mind
After=network.target

[Service]
Type=simple
ExecStart=/home/ivanh/.local/bin/llama-server \
    -m /home/ivanh/models/nomic-embed-text-v1.5.Q4_K_M.gguf \
    --embeddings \
    --pooling mean \
    -c 2048 \
    -b 2048 \
    -ub 2048 \
    --host 127.0.0.1 \
    --port 8089
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
